# debian/sfha.service
[Unit]
Description=sfha - Système de Haute Disponibilité léger
Documentation=https://github.com/serverflow/sfha
After=network.target corosync.service
Requires=corosync.service
ConditionPathExists=/etc/sfha/config.yml

[Service]
Type=simple
ExecStart=/usr/bin/sfha run
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sfha

# Désactivation de TOUTES les restrictions namespace pour compatibilité LXC
PrivateTmp=no
ProtectHome=no
ProtectSystem=no
PrivateDevices=no
ProtectHostname=no
ProtectClock=no
ProtectKernelTunables=no
ProtectKernelModules=no
ProtectKernelLogs=no
ProtectControlGroups=no
RestrictNamespaces=no

[Install]
WantedBy=multi-user.target
